# quick-start/fast-dashboard/fast-dashboard.container
[Unit]
Description=Fast Dashboard
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/shrekbytes/fast-dashboard:latest
ContainerName=fast-dashboard
AutoUpdate=registry

# change ip and port from config.yml
Network=host

# Environment file (create .env in the folder if you need secrets in config)
EnvironmentFile=%h/fast-dashboard/.env

# Persistent storage (same layout as compose: one folder with config/ and optional assets/)
Volume=%h/fast-dashboard/config:/app/config:Z
Volume=%h/fast-dashboard/assets:/app/assets:Z
Volume=/etc/localtime:/etc/localtime:ro

# Port mapping (will not work when on Network=host)
# PublishPort=9999:8080

[Service]
# Restart behavior
Restart=on-failure


[Install]
WantedBy=default.target
